# recipes-core/battery-service/files/librescoot-battery.service
[Unit]
Description=Librescoot Battery Service
Requires=redis.service

[Service]
Type=idle
User=root
Group=root
ExecStart=/usr/bin/battery-service
Restart=always
KillMode=mixed

[Install]
WantedBy=multi-user.target
